# debian/postrm
#! /bin/sh

# 1. Remove system startup links
update-rc.d  -f gotun remove
exit 0
